# systemd/oscillatord.service.in
[Unit]
Description=Daemon responsible of disciplining an oscillator based on a 1pps phase error device.

[Service]
Type=simple
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/oscillatord /etc/oscillatord.conf

[Install]
WantedBy=multi-user.target
